# deployment/config/gunicorn.conf
; =====================
;  manage.py supervisor
; =====================

[program:gunicorn]
command=/home/ohc/.virtualenvs/odonto/bin/gunicorn odonto.wsgi -c /usr/lib/ohc/etc/gunicorn_conf.py
directory=%(here)s/../../
;user=nobody
numprocs=1
stdout_logfile=/usr/lib/ohc/log/gunicorn.log
stderr_logfile=/usr/lib/ohc/log/gunicorn.log
autostart=true
autorestart=true
startsecs=10
loglevel=DEBUG
